# project/config/production/crontab/api_frame
* * * * * root /bin/bash /var/www/api_frame/project/tool/check_update.sh
