# project/config/production/supervisor/caddy.conf
[program:caddy]
command=/usr/local/bin/caddy -conf caddy_file
directory=/root
restart=always
user=root
