# project/config/development/supervisor/queue_job_watch.conf
[program:queue_job_watch]
command=/bin/sh /var/www/api_frame/project/tool/development/queue_job_watch_by_md5.sh
stopwaitsecs=5
restart=always
user=root
